# scripts/co.9arm.kuro-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>co.9arm.kuro-publish</string>

	<key>ProgramArguments</key>
	<array>
		<string>/Users/tpatinya/kuro/scripts/publish-runner.sh</string>
	</array>

	<!-- Saturday 00:00 (Weekday: 0=Sun, 1=Mon, ..., 6=Sat). -->
	<key>StartCalendarInterval</key>
	<dict>
		<key>Weekday</key>
		<integer>6</integer>
		<key>Hour</key>
		<integer>0</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>

	<key>StandardOutPath</key>
	<string>/Users/tpatinya/Library/Logs/kuro-publish.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/tpatinya/Library/Logs/kuro-publish.err.log</string>
</dict>
</plist>
